# files/system/etc/systemd/system/system-flatpak-setup.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target NetworkManager-wait-online.service

[Service]
Restart=on-failure
RestartSec=30
